# supervisord.conf
[supervisord]
nodaemon=true

[program:syslogd]
directory=/srv
command=npm syslogd --production

[program:amqpd]
directory=/srv
command=npm amqpd --production
