# supervisord.conf
[supervisord]
nodaemon=true

[program:syslogd]
directory=/srv
command=npm run-script syslogd --production
